# debian/maemo-system-services.welcome.init
#!/sbin/openrc-run

start() {
	einfo "Welcome to Maemo Leste"
}
