# debian/maemo-system-services.welcome.service
[Service]
ExecStart=/bin/echo -e "\nWelcome to your systemd powered internet tablet."
StandardOutput=tty
