# debian/maemo-system-services.welcome.init
#!/sbin/openrc-run
#### BEGIN INIT INFO
# Provides:          welcome
# Required-Start:    $remote_fs
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description:
### END INIT INFO

start() {
	einfo "\nWelcome to your upstart powered internet tablet.\n"
}
