# debian/maemo-system-services.welcome.init
#!/sbin/openrc-run

start() {
	einfo "\nWelcome to your upstart powered internet tablet.\n"
}
